# installer/systemd/wcp360.service
# ==============================================================================
# WCP360 Systemd Service Unit
# ------------------------------------------------------------------------------
# Runs WCP360 daemon as isolated system user.
#
# Security:
#   - Runs under wcp user
#   - Automatic restart
#   - Network dependency enforced
# ==============================================================================
[Unit]
Description=WCP360 Daemon
After=network.target

[Service]
User=wcp
Group=wcp
ExecStart=/opt/wcp360/bin/wcpd
Restart=always

[Install]
WantedBy=multi-user.target
